# install/videobridge.service
[Unit]
Wants=network.target
[Service]
ExecStart=/opt/videobridge/run-videobridge.sh
User=nobody
Group=plugdev
[Install]
WantedBy=multi-user.target
